# net-dns/inadyn-mt/files/inadyn-mt.initd
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/inadyn-mt"
command_args="${OPTS} ${CONFIG}"
description="Dynamic DNS client"
name="Dynamic DNS client"
user="${INADYN_USER:-inadyn}"
group="${INADYN_GROUP:-inadyn}"
pidfile="/run/inadyn-mt/${RC_SVCNAME}.pid"

depend() {
	need net
}
