# net-misc/ssmd/files/ssmd.initd
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/sbin/ssmd"
command_args="${OPTS} ${CONFIG}"
description="SNMP Switch Management Daemon"
name="SNMP Switch Management Daemon"
user="${SSMD_USER:-ssmd}"
group="${SSMD_GROUP:-ssmd}"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	need net
}
