# app-metrics/aws-otel-collector/files/aws-otel-collector.initd
#!/sbin/openrc-run
# Copyright 1999-2022 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/opt/aws/aws-otel-collector/bin/aws-otel-collector"
command_args="${AOC_OPTS}"
command_background="yes"
description="Amazon AWS Opentelemetry Collector"
name="Amazon AWS Opentelemetry Collector"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	need net
}
