# net-analyzer/netgauge/files/netgauge.initd
#!/sbin/openrc-run
# Copyright 1999-2022 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/opt/netgauge/OoklaServer"
pidfile="/run/${RC_SVCNAME}.pid"
command_args="${OPTS} --pidfile=$pidfile"
name="Ookla NetGauge"
retry=20

depend() {
	need net
}
