# net-vpn/awsvpnclient/files/awsvpnclient.initd
#!/sbin/openrc-run
# Copyright 1999-2025 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/opt/awsvpnclient/Service/ACVC.GTK.Service"
command_background="yes"
name="AWS VPN Clinet"
description="AWS VPN Clinet"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	need net
}
