# dev-db/litestream/files/litestream.initd
#!/sbin/openrc-run
# Copyright 1999-2022 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/litestream"
name="Litestream"
description="Litestream"
command_args="${LS_OPTS}"

depend() {
	need net
}
