# x11-misc/logkeys/files/logkeys.service
[Unit]
Description=Linux keylogger

[Service]
Type=forking
GuessMainPID=yes
ExecStop=/usr/bin/logkeys --kill

[Install]
WantedBy=multi-user.target
