# x11-misc/logkeys/files/logkeys.initd
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/logkeys"
command_args="-s ${ARGS}"
command_background="true"
description="Linux keylogger"
name="Linux keylogger"
pidfile="/run/${RC_SVCNAME}.pid"
