# net-misc/msd/files/msd.initd
#!/sbin/openrc-run
# Copyright 1999-2022 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/msd"
command_args="${MSD_OPTS}"
description="Multi Stream daemon"
name="Multi Stream daemon"
user="${MSD_USER:-msd}"
group="${MSD_GROUP:-msd}"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	need net
}
