# net-dns/r53ddns/files/r53ddns.initd
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command=/usr/bin/r53ddns
command_args="${OPTS}"
description="Amazon Route 53 updater"
name="Amazon Route 53 updater"
pidfile=/run/r53ddns.pid
start_stop_daemon_args="-q"

depend() {
	need net
}
