# app-backup/urbackup-client/files/urbackup-client.initd
#!/sbin/openrc-run
# Copyright 1999-2022 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

pidfile=/run/urbackup-client.pid
command=/usr/sbin/urbackupclientbackend
command_args="--config ${URBACKUP_CONFIG} --pidfile ${pidfile} --daemon ${URBACKUP_OPTS}"
description="UrBackup client"
name="UrBackup client"

depend() {
	need net
}
