# net-vpn/innernet/files/innernet-client.initd
#!/sbin/openrc-run
# Copyright 1999-2025 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/innernet"
name="Innernet client"
description="Innernet client"
command_args="up ${INNERNET_OPTS}"

depend() {
	need net
}
