# net-vpn/innernet/files/innernet-server.initd
#!/sbin/openrc-run
# Copyright 1999-2025 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/innernet-server"
name="Innernet server"
description="Innernet server"
command_args="serve ${INNERNET_OPTS}"

depend() {
	need net
}
