# net-misc/ntripserver/files/ntripserver.initd
#!/sbin/openrc-run
# Copyright 1999-2022 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/ntripserver"
command_args="${NTRIP_OPTS} -m ${MOUNTPOINT} ${CONFIG}"
command_background="yes"
description="Ntrip Command Line Server"
name="Ntrip Command Line Server (${MOUNTPOINT})"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	need net
}
