# net-misc/mcproxy/files/mcproxy.initd
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/mcproxy"
command_args="-f /etc/mcproxy.conf ${OPTS}"
command_background="yes"
description="Start mcproxy"
pidfile="/run/mcproxy.pid"
start_stop_daemon_args="-q"

depend() {
	need net
}
