# net-misc/mcproxy/files/mcproxy.service
[Unit]
Description=mcproxy
After=network.target

[Service]
ExecStart=/usr/sbin/mcproxy /etc/mcproxy.conf

[Install]
WantedBy=multi-user.target
