# net-misc/inverting-proxy/files/inverting-proxy-agent.initd
#!/sbin/openrc-run
# Copyright 1999-2022 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/proxy-forwarding-agent"
command_args="${INVP_OPTS}"
command_background="yes"
command_user="${INVP_USER:-root}:${INVP_GROUP:-root}"
description="Inverting Proxy agent"
name="Inverting Proxy agent"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	need net
}
