# www-apps/noc-tower/files/noc-tower.initd
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/tower-web"
command_args="${OPTS}"
command_background="yes"
description="NOC Tower"
group="${NOC_GROUP:-noc}"
name="NOC Tower"
pidfile="/run/noc/${RC_SVCNAME}.pid"
user="${NOC_USER:-noc}"

depend() {
	need net
}
