# net-misc/fiche/files/fiche.initd
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/fiche"
command_args="${OPTS}"
command_backround="yes"
description="Pastebin server"
name="Pastebin server"
user="${FICHE_USER:-fiche}"
group="${FICHE_GROUP:-fiche}"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	need net
}
